# hub/claude-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.claude-monitor</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/user/.nvm/versions/node/v24.13.1/bin/npm</string>
    <string>start</string>
    <string>--prefix</string>
    <string>/Users/user/Documents/Development/claude-monitor/hub</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/user/Documents/Development/claude-monitor</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/user/Library/Logs/claude-monitor.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/user/Library/Logs/claude-monitor.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/user/.nvm/versions/node/v24.13.1/bin:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
  </dict>
</dict>
</plist>
